# ubuntu-deploy-complete/elontalk.service
[Unit]
Description=ELonTalk Customer Service System
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=ubuntu
ExecStart=/opt/elontalk/customer-service-backend
WorkingDirectory=/opt/elontalk
Environment=RUST_LOG=info
Environment=DATABASE_URL=sqlite:customer_service.db
EnvironmentFile=-/opt/elontalk/.env

# 安全设置
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/elontalk

# 资源限制
LimitNOFILE=65535
LimitNPROC=4096

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=elontalk

[Install]
WantedBy=multi-user.target